# private_dot_config/systemd/user/wayland-wm@hyprland.desktop.service.d/logging.conf.tmpl
[Service]
LogsDirectory=hyprland
StandardOutput=append:%h/.local/state/hyprland/hyprland.log
StandardError=append:%h/.local/state/hyprland/hyprland.log
