# private_dot_config/systemd/system/backup.timer.tmpl
[Unit]
Description=Schedule Borg backups

[Timer]
OnCalendar={{ .backup.on_calendar }}
Persistent=true
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
